# gunicorn.service
# CONFIGURATION
# PROJECT_NAME: chatapp
# PROJECT_DIR: /var/www/chatapp
# WSGI_MODULE: chatapp.wsgi:application

[Unit]
Description=Gunicorn daemon for chatapp Django application
After=network.target

[Service]
User=www-data
Group=www-data
UMask=0007
WorkingDirectory=/var/www/chatapp
Environment="PATH=/var/www/chatapp/venv/bin"
EnvironmentFile=/var/www/chatapp/.env
ExecStart=/var/www/chatapp/venv/bin/gunicorn \
    --workers 3 \
    --bind unix:/var/www/chatapp/gunicorn.sock \
    --timeout 60 \
    --access-logfile /var/www/chatapp/logs/gunicorn-access.log \
    --error-logfile /var/www/chatapp/logs/gunicorn-error.log \
    chatapp.wsgi:application

[Install]
WantedBy=multi-user.target
